# recipes-custom/busybox/files/udhcpc.service
[Unit]
Description=DHCP Client

[Service]
Type=simple
ExecStart=/usr/sbin/udhcpc

[Install]
WantedBy=multi-user.target
